# config/systemd/user/checkupdates.timer
[Timer]
OnBootSec=60min
OnUnitActiveSec=60min

[Install]
WantedBy=default.target
