# config/systemd/user/syncmail.service
[Unit]
Description="Mail syncing script"

[Service]
Type=oneshot
ExecStart=/usr/bin/mailsync
User=pascal
Group=pascal
Environment=DISPLAY=:0
