# config/systemd/user/syncmail.service
[Unit]
Description="Mail syncing script"

[Service]
Type=oneshot
ExecStart=/usr/local/bin/syncmail

