# config/systemd/user/syncmail.timer
[Unit]
Description=Run syncmail every five minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=sway-session.target
